# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:backend]
command=/app/probe-server
directory=/app
autostart=true
autorestart=true
startsecs=5
; Cap restart attempts so supervisord surfaces a real failure (e.g. DB
; corruption) instead of hot-looping the backend forever and flooding logs.
startretries=5
; Forward docker stop / compose down as SIGTERM and give the backend up to
; 30 s to flush the bbolt DB and close the file cleanly. If we let supervisord
; send SIGKILL too early, an in-flight db.Update can leave metrics.db in a
; state that triggers the old bbolt "invalid freelist page" failure mode on
; the next boot.
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT="8080"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=5
startretries=5
stopsignal=QUIT
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

